# meta-mcst/meta-common/recipes-mcst/active-gpio-on/files/active-gpio-on.service
[Unit]
Description=Turn on ACTIVE# and turn off BMC_BOOT# LEDs
After=getty.target

[Service]
ExecStart=/libexec/active-gpio-on start
ExecStop=/libexec/active-gpio-on stop
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
