# meta-mcst/meta-common/recipes-poky/openssh/files/sshd.socket
[Unit]
Conflicts=sshd.service
Wants=sshdgenkeys.service

[Socket]
ExecStartPre=@BASE_BINDIR@/mkdir -p /var/run/sshd
ListenStream=22
ListenStream=2022
Accept=yes

[Install]
WantedBy=sockets.target
