# meta-mcst/meta-common/recipes-mcst/rtcread/files/rtcread.service
[Unit]
Description=Read RTC from motherboard
After=sysinit.target

[Service]
Type=oneshot
ExecStart=/libexec/rtcread

[Install]
WantedBy=basic.target
