# meta-mcst/meta-common/recipes-mcst/tmux-logger/files/logrotate-sol.timer
[Unit]
Description=Rotate SoL log each minute
Requires=logrotate-sol.service

[Timer]
Unit=logrotate-sol.service
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target
